# slash/etc/s6-overlay/s6-rc.d/php-fpm83/run
#!/command/with-contenv sh
exec /usr/sbin/php-fpm83 --nodaemonize --fpm-config /etc/php83/php-fpm.conf
