# slash/etc/s6-overlay/s6-rc.d/php-fpm84/run
#!/command/with-contenv sh
exec /usr/sbin/php-fpm84 --nodaemonize --fpm-config /etc/php84/php-fpm.conf
