# slash/etc/s6-overlay/s6-rc.d/angie/run
#!/command/with-contenv sh
exec /usr/sbin/angie -c /etc/angie/angie.conf
